# deploy/systemd/portfolio-backup.service
[Unit]
Description=Portfolio backup (MariaDB + n8n Postgres + WP uploads + configs)
After=docker.service

[Service]
Type=oneshot
ExecStart=/opt/portfolio/scripts/backup.sh
