# deploy/systemd/portfolio-restore-verify.service
[Unit]
Description=Portfolio backup restore verification (sandbox)
After=docker.service

[Service]
Type=oneshot
ExecStart=/opt/portfolio/scripts/restore-test.sh
